# conf/clean_blobrestore_jobs.cron
@weekly root find /home/blobrestore/processed_jobs/ -mindepth 1 -maxdepth 1 -type d -mtime +7 -exec rm -rf {} \;
