# conf/membase-backup-incr.cron
00,20,40 * * * * /opt/membase/membase-backup/mbbackup-incr > /dev/null 2>&1
